# homeassistant-mcp-http/rootfs/etc/s6-overlay/s6-rc.d/mcp-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant MCP HTTP Server
# Runs the MCP server
# ==============================================================================

bashio::log.info "Starting MCP HTTP Server..."

# Get OAuth URL from addon config
OAUTH_CLIENT_URL=$(bashio::config 'oauth_client_url')

if bashio::var.is_empty "${OAUTH_CLIENT_URL}"; then
    bashio::log.fatal "oauth_client_url not configured!"
    bashio::log.fatal "Please configure your DuckDNS URL in addon options"
    bashio::exit.nok
fi

# Export environment variables
export OAUTH_CLIENT_URL="${OAUTH_CLIENT_URL}"
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
export INGRESS_PATH="${INGRESS_PATH}"

bashio::log.info "OAuth Client URL: ${OAUTH_CLIENT_URL}"
bashio::log.info "Ingress Path: ${INGRESS_PATH}"

# Start server
cd /app || bashio::exit.nok "Cannot change to /app directory"
exec node dist/http-server.js
